# config/systemd/user/kmonad.service
[Unit]
Description=Kmonad

[Service]
ExecStart=kmonad /home/huy/.config/kmonad/config.kbd
Restart=on-failure

[Install]
WantedBy=default.target
