# huymai/services/n8n.service
[Unit]
Description=n8n Workflow Automation
Documentation=https://docs.n8n.io
After=network.target

[Service]
ExecStart=/usr/bin/npx n8n
Restart=always
RestartSec=10
Environment=NODE_ENV=production
WorkingDirectory=/path/to/n8n/project
User=n8n
Group=n8n

# Logging
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=n8n

# Limit resources
LimitNOFILE=10000

[Install]
WantedBy=multi-user.target
