# systemd_services/ibus.service
[Unit]
Description=Ibus

[Service]
ExecStart=ibus-daemon
Restart=on-failure

[Install]
WantedBy=default.target
